# system/files/systemd/dconf.timer
[Unit]
Description=Run dconf backup script

[Timer]
OnCalendar=daily
OnActiveSec=2
Persistent=true

[Install]
WantedBy=backups.target
