# system/files/systemd/dconf.timer
[Unit]
Description=Run dconf backup script

[Timer]
OnCalendar=daily
Persistent=true
Unit=dconf.service

[Install]
WantedBy=backup.target
