# system/files/systemd/dconf.timer
[Unit]
Description=Run dconf backup script

[Timer]
OnStartupSec=10min
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
