# system/files/systemd/duplicity.timer
[Unit]
Description=Run duplicity backup script
After=dconf.service

[Timer]
OnCalendar=daily
OnActiveSec=5
Persistent=true

[Install]
WantedBy=backups.target
